# pdu-webhook.service
[Unit]
Description=PDU Webhook Handler
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/PDU-NEW
Environment=PATH=/opt/PDU-NEW/pdu_env/bin
ExecStart=/opt/PDU-NEW/pdu_env/bin/python3 webhook_handler.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
